# motion/.footprint
drwxr-xr-x	root/root	etc/
drwxr-xr-x	root/root	etc/motion/
-rw-r-----	root/motion	etc/motion/motion.conf
-rw-r-----	root/motion	etc/motion/thread1.conf
-rw-r-----	root/motion	etc/motion/thread2.conf
-rw-r-----	root/motion	etc/motion/thread3.conf
-rw-r-----	root/motion	etc/motion/thread4.conf
drwxr-xr-x	root/root	etc/rc.d/
-rwxr-xr-x	root/root	etc/rc.d/motion
drwxr-xr-x	root/root	usr/
drwxr-xr-x	root/root	usr/bin/
-rwxr-xr-x	root/root	usr/bin/motion
drwxr-xr-x	root/root	usr/man/
drwxr-xr-x	root/root	usr/man/man1/
-rwxr-xr-x	root/root	usr/man/man1/motion.1.gz
drwxr-xr-x	root/root	var/
drwxr-xr-x	root/root	var/run/
drwxr-xr-x	motion/motion	var/run/motion/
